# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-rescue.ubifs.extracted/squashfs-root/usr/lib/opkg/info/kmod-broadcom-voicedriver-4.16L.list
/etc/telephonyProfiles.d/MNE_profile.xml
/etc/telephonyProfiles.d/USA_profile.xml
/etc/telephonyProfiles.d/HUN_profile.xml
/etc/telephonyProfiles.d/FRA_profile.xml
/etc/telephonyProfiles.d/TWN_profile.xml
/etc/modules.d/93-broadcom-voicedriver-4.16L
/etc/telephonyProfiles.d/BEL_profile.xml
/etc/telephonyProfiles.d/SVK_profile.xml
/etc/telephonyProfiles.d/CHE_profile.xml
/etc/telephonyProfiles.d/ITA_profile.xml
/etc/telephonyProfiles.d/ETS_profile.xml
/etc/telephonyProfiles.d/UAE_profile.xml
/etc/telephonyProfiles.d/AUT_profile.xml
/etc/telephonyProfiles.d/CYP_profile.xml
/etc/telephonyProfiles.d/GRC_profile.xml
/etc/telephonyProfiles.d/NLD_profile.xml
/etc/telephonyProfiles.d/SVN_profile.xml
/etc/telephonyProfiles.d/FIN_profile.xml
/etc/telephonyProfiles.d/SWE_profile.xml
/etc/telephonyProfiles.d/BRA_profile.xml
/etc/telephonyProfiles.d/DNK_profile.xml
/etc/telephonyProfiles.d/POL_profile.xml
/etc/telephonyProfiles.d/GBR_profile.xml
/etc/telephonyProfiles.d/MKD_profile.xml
/etc/telephonyProfiles.d/DEU_profile.xml
/etc/telephonyProfiles.d/MEX_profile.xml
/etc/telephonyProfiles.d/IRL_profile.xml
/etc/telephonyProfiles.d/CHN_profile.xml
/etc/telephonyProfiles.d/PRT_profile.xml
/etc/telephonyProfiles.d/ROU_profile.xml
/etc/telephonyProfiles.d/T57_profile.xml
/etc/telephonyProfiles.d/CZH_profile.xml
/lib/modules/3.4.11-rt19/endpointdd.ko
/etc/telephonyProfiles.d/ESP_profile.xml
/etc/telephonyProfiles.d/NOR_profile.xml
/etc/telephonyProfiles.d/CHL_profile.xml
/etc/init.d/voice
/etc/telephonyProfiles.d/MLT_profile.xml
/etc/telephonyProfiles.d/HRV_profile.xml
/etc/telephonyProfiles.d/AUS_profile.xml
/etc/telephonyProfiles.d/IND_profile.xml
/etc/telephonyProfiles.d/JPN_profile.xml
/etc/telephonyProfiles.d/SGP_profile.xml
/etc/telephonyProfiles.d/NZL_profile.xml
