# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-rescue.ubifs.extracted/squashfs-root/usr/lib/opkg/info/lighttpd.conffiles
/etc/lighttpd/lighttpd.conf
/etc/lighttpd/conf.d/mimetype.conf
/etc/lighttpd/conf.d/bind.conf
/etc/lighttpd/conf.d/wan-ip-addr.conf
/etc/lighttpd/conf.d/redirection.conf
/etc/lighttpd/conf.d/rest.conf
/etc/init.d/lighttpd
/etc/lighttpd/ssl/remote/self-signed.pem
