# ubi/ubifs-root/firware.bin/_img-0_vol-operational.ubifs.extracted/squashfs-root/etc/init.d/network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=11

start() {
	/sbin/ifconfig lo 127.0.0.1 netmask 255.0.0.0 broadcast 127.255.255.255 up
}

