# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/device-info-smu.list
/etc/init.d/syslogfirewall
/usr/lib/xmo/plugins/libdevice-info-smu.so
/etc/init.d/syslog
