# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/dhcp-server.list
/usr/bin/omshell
/etc/init.d/isc-dhcp
/usr/sbin/dhcpd
