# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/gui-core.list
/usr/share/gui-core/www/gui/css/images/icons-36-white.png
/usr/share/gui-core/www/gui/js/jquery.js
/usr/share/gui-core/www/gui/js/jquery.form.js
/usr/share/gui-core/www/gui/css/jquery.jqGrid.css
/usr/share/gui-core/www/gui/js/gui-api.js
/usr/share/gui-core/www/gui/js/jquery.metadata.js
/usr/share/gui-core/www/gui/js/jquery.random.js
/etc/lighttpd/conf.d/gui-version.conf
/usr/share/gui-core/www/gui/css/images/ajax-loader.gif
/usr/share/gui-core/www/gui/js/jquery-ui-1.7.2.custom.js
/etc/init.d/gui-firmware
/usr/share/gui-core/www/gui/css/jquery-ui.css
/usr/share/gui-core/www/gui/js/jquery.cookie.js
/usr/share/gui-core/www/gui/css/images/icons-18-black.png
/usr/share/gui-core/www/cgi/xml-api
/usr/share/gui-core/www/gui/css/images/icons-36-black.png
/usr/share/gui-core/www/gui/css/jquery.mobile.css
/usr/share/gui-core/www/gui/js/gui-core.js
/usr/share/gui-core/www/gui/js/jquery.mobile.js
/usr/share/gui-core/www/gui/js/jquery-ui-1.7.3.custom-sc.js
/usr/share/gui-core/www/gui/js/jquery-1.3.2-sc.js
/usr/share/gui-core/www/gui/js/jquery-1.7.2-sc.js
/usr/share/gui-core/www/gui/js/gui-widgets.js
/usr/share/gui-core/www/gui/js/jquery-utils.js
/usr/share/gui-core/www/gui/js/jquery.jcarousel.js
/usr/share/gui-core/www/gui/js/gui-common.js
/usr/share/gui-core/www/gui/js/jquery-ui.js
/usr/share/gui-core/www/gui/js/grid.locale-en.js
/usr/share/gui-core/www/gui/msg/fr/gui-api-msg.json
/usr/share/gui-core/www/gui/css/jquery-ui-1.8.20.custom.css
/usr/share/gui-core/www/cgi/upload
/usr/share/gui-core/www/cgi/json-req
/usr/share/gui-core/www/gui/js/jquery.jqGrid.js
/usr/share/gui-core/www/gui/css/images/ajax-loader.png
/usr/share/gui-core/www/gui/js/jquery.printf.js
/etc/gui-core.conf
/usr/share/gui-core/www/gui/js/md5.js
/usr/share/gui-core/www/gui/js/jquery-ui-1.8.20.custom.js
/usr/share/gui-core/www/gui/msg/en/gui-api-msg.json
/usr/share/gui-core/www/gui/msg/nl/gui-api-msg.json
/usr/share/gui-core/www/gui/js/json2.js
/usr/share/gui-core/www/gui/css/gui-common.css
/usr/share/gui-core/www/gui/css/images/icons-18-white.png
/usr/share/gui-core/www/gui/css/gui-core.css
/usr/share/gui-core/www/gui/js/jquery.tmpl.js
