# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/hg6d.list
/usr/bin/hg6d_start
/usr/bin/hg6d
/etc/init.d/hg6d
