# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/lighttpd.list
/usr/lib/lighttpd/mod_indexfile.so
/etc/lighttpd/conf.d/rest.conf
/etc/lighttpd/conf.d/wan-ip-addr.conf
/etc/init.d/lighttpd
/etc/lighttpd/conf.d/redirection.conf
/etc/lighttpd/ssl/remote/self-signed.pem
/usr/sbin/lighttpd
/usr/lib/lighttpd/mod_staticfile.so
/etc/lighttpd/conf.d/bind.conf
/etc/lighttpd/conf.d/mimetype.conf
/usr/sbin/lighttpd-angel
/lib/upgrade/keep.d/lighttpd
/usr/lib/lighttpd/mod_dirlisting.so
