# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/usb-smu.list
/usr/bin/usb2smu
/etc/init.d/usb_mod
/usr/lib/xmo/plugins/libusb-smu.so
/etc/init.d/sata_mod
