# ubi/ubifs-root/firmware_bin.extracted/_img-0_vol-operational.ubifs.extracted/squashfs-root/etc/init.d/isc-dhcp
#!/bin/sh /etc/rc.common
START=60

start() {
	chmod 666 /dev/null
}

stop() {
	echo stop device /dev/null
}
