# ubi/ubifs-root/firware.bin/_img-0_vol-operational.ubifs.extracted/squashfs-root/usr/lib/opkg/info/firewall-utils.list
/usr/lib/libfirewall.so
/etc/init.d/firewall-utils
/usr/lib/libfirewall.so.1.0.0
/usr/sbin/firewall-app
/usr/lib/libfirewall.so.1
/etc/init.d/firewall-app
